# rsteye/etc/systemd/system/rsteyeapp.service
[Unit]
Description=RstEyeApp Daemon
After=network.target

[Service]
ExecStart=/usr/bin/RstEyeApp
Restart=always

[Install]
WantedBy=multi-user.target