# debian_pkg_files/postinst
#!/bin/bash
set -e
systemctl enable rsteye.service
systemctl start rsteye.service
