# debian_pkg_files/postrm
#!/bin/bash
set -e
systemctl stop rsteye.service
systemctl disable rsteye.service
