# rsteye/DEBIAN/prerm
#!/bin/bash

set -e

# Stop and disable the service
systemctl stop rsteyeapp.service
systemctl disable rsteyeapp.service

exit 0